# server_config/gun_supervisor.conf
[program:gunicorn]
command=gunicorn run_server:application
directory: /home/ubuntu/note-tagger_package/note_tagger
user=nobody
autostart=true
autorestart=true
redirect_stderr=true