# microblog_package/server_config/gun_supervisor.conf
[program:gunicorn]
command=gunicorn -b 127.0.0.1:5000 run_server:app
directory: /home/ubuntu/microblog_package
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=true